# contrib/systemd/towel-worker.service
[Unit]
Description=Towel AI Worker — auto-discovers GPU, models, and controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/towel worker
Restart=always
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# GPU access
SupplementaryGroups=video render

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=towel-worker

[Install]
WantedBy=default.target
